# src/vault/build/modules/build/files/etc/supervisor/conf.d/vault.conf
[program:vault]
command = /usr/local/bin/vault server -config=/etc/vault/vault.hcl
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
redirect_stderr = true
